# system_files/etc/profile.d/rpm-ostree-search-hl.sh
alias rpm-ostree="rpm-ostree-search-hl"
